# fastboot/systemd/system/data.mount
[Unit]
Description=Fastboot /data tmpfs

[Mount]
What=tmpfs
Where=/data
Type=tmpfs
Options=X-mount.mkdir
